# deploy/systemd/okx-remote-auto-update.service
[Unit]
Description=Auto update OKX Remote Node from GitHub
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/okx-local-app
Environment=OKX_REMOTE_APP_DIR=/opt/okx-local-app
Environment=OKX_REMOTE_BRANCH=main
Environment=OKX_REMOTE_GIT_REMOTE=origin
Environment=OKX_REMOTE_SERVICE_NAME=okx-remote-node
ExecStart=/opt/okx-local-app/deploy/auto-update-remote-node.sh

