# deploy/systemd/okx-remote-node.service
[Unit]
Description=OKX Local App Remote Execution Node
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/okx-local-app
EnvironmentFile=-/etc/okx-remote-node.env
Environment=OKX_LOCAL_APP_HOST=127.0.0.1
Environment=OKX_LOCAL_APP_PORT=18765
Environment=OKX_LOCAL_APP_DATA_DIR=/opt/okx-local-app/data-remote
ExecStart=/opt/okx-local-app/scripts/start-remote-node.sh
Restart=always
RestartSec=3
User=ecs-assist-user
Group=ecs-assist-user

[Install]
WantedBy=multi-user.target
